# pkgbuild/auto-rules-io.install
post_install() {
    systemctl enable big-btrfs-balance.service
}

post_upgrade() {
    systemctl enable big-btrfs-balance.service
    sed -i 's|,space_cache=v2,autodefrag,compress-force=zstd:2,commit=15(\|,discard=async,ssd)|,compress-force|g' /etc/fstab 
}

pre_remove() {
    systemctl disable big-btrfs-balance.service
} 
